# packages/sx-systemd/skel/lib/systemd/user/sx-session.target
# SPDX-FileCopyrightText: 2025 Daniel Sampliner <samplinerD@gmail.com>
#
# SPDX-License-Identifier: AGPL-3.0-or-later

[Unit]
Description=sx session

DefaultDependencies=no
BindsTo=graphical-session.target
Before=graphical-session.target
